# com.flightwatcher.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS LaunchAgent: runs the health check at 6:15pm IST (Mon-Fri) to verify the
  6pm flight scan completed and sent the Telegram report. If it failed,
  sends an alert Telegram message.

  This is the reference copy. The live copy lives at:
    ~/Library/LaunchAgents/com.flightwatcher.health.plist

  Install / reload:
    cp com.flightwatcher.health.plist ~/Library/LaunchAgents/
    launchctl unload ~/Library/LaunchAgents/com.flightwatcher.health.plist 2>/dev/null
    launchctl load -w ~/Library/LaunchAgents/com.flightwatcher.health.plist

  Run once now to test:
    launchctl start com.flightwatcher.health
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flightwatcher.health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/sarathkrishnan/flight-watcher/health_check.py</string>
        <string>--alert</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sarathkrishnan/flight-watcher</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 18:15 IST, Mon-Fri — verify 6pm report was sent, alert if not -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/sarathkrishnan/flight-watcher/logs/health-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sarathkrishnan/flight-watcher/logs/health-check.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
